# filas.conf
[program:projeto]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/filas/artisan queue:work --daemon
autostart=true
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/html/filas/storage/logger.log